# deploy/configs/pufferfish.target
[Unit]
Description=Pufferfish Target
Requires=graphical.target
After=graphical.target